# backend/parkingspots-api.service
[Unit]
Description=ParkingSpots API Service (Multi-Worker)
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=dalas
Group=dalas
WorkingDirectory=/home/dalas/ParkingSpots/backend
Environment="ENABLE_BACKGROUND_TASKS=false"
Environment="DATABASE_URL=postgresql+asyncpg://parking_user:parking_secure_2026@localhost:5432/parkingspots"
ExecStart=/home/dalas/ParkingSpots/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 12
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/dalas/ParkingSpots/backend

[Install]
WantedBy=multi-user.target
